# systemd/.config/systemd/user/offlineimap.service
[Unit]
Description=offlineimap sync daemon
After=network.target

[Service]
ExecStart=/usr/bin/offlineimap -u quiet
KillSignal=SIGUSR2
Restart=always
